# llama-switch.service
[Unit]
Description=Llama Switch Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/AI/LLama Switch
ExecStart=/usr/bin/python3 server.py --config config.yaml --watch
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
